# roles/rethinkdb_client/templates/etc/cron.daily/unassign_stf_devices
#!/usr/bin/env bash

recli -d stf 'r.table("devices").filter(r.row("owner").ne(null)).update({owner: null})'
